# systemd/mobilealerts.service
[Unit]
Description=Mobile Alerts Proxy Server
After=network.target

[Service]
# WICHTIG: Ersetze /usr/bin/node durch den tats채chlichen Pfad, den 'which node' liefert.
# WICHTIG: F체ge alle deine Kommandozeilen-Argumente (IP, Port) hinzu.
ExecStart=/usr/bin/node /home/raspiroman/project/MMMMobileAlerts/maserver/mobilealerts.js

# Das Verzeichnis, in dem sich die Skript-Dateien (wie package.json, conf) befinden.
WorkingDirectory=/home/raspiroman/project/MMMMobileAlerts/maserver/

# Der Benutzer, unter dem der Dienst ausgef체hrt werden soll (dein regul채rer Benutzer).
User=raspiroman

# Den Dienst bei Absturz automatisch neu starten.
Restart=always

# Warte 5 Sekunden nach einem Fehler, bevor neu gestartet wird.
RestartSec=30

# Standardausgabe (console.log) an das Systemd Journal senden.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target